# scripts/service/deluged/run
#!/bin/bash

sv start openvpn || exit 1
sleep 5

exec /sbin/setuser torrents /usr/bin/deluged -c /config/deluge -d --loglevel=info
